# growberry-wifi-watchdog.timer
[Unit]
Description=Run Growberry WiFi Watchdog every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
